# kiwix/rootfs/etc/services.d/nginx.disabled/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Kiwix
# Runs nginx for ingress support
# ==============================================================================

# Create all directories nginx might need in /tmp (always writable)
mkdir -p /tmp/nginx/client_body \
         /tmp/nginx/proxy \
         /tmp/nginx/fastcgi \
         /tmp/nginx/uwsgi \
         /tmp/nginx/scgi \
         /run/nginx

# Wait for Kiwix to be available
while ! nc -z 127.0.0.1 8090; do
    bashio::log.info "Waiting for Kiwix to start..."
    sleep 2
done

bashio::log.info "Starting nginx for ingress support..."

# Start nginx directly (let it handle temp directories)
exec nginx
