# kiwix/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start nginx for ingress support
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Create temp directories for nginx
mkdir -p /tmp/nginx/client_body
mkdir -p /tmp/nginx/proxy
mkdir -p /tmp/nginx/fastcgi
mkdir -p /tmp/nginx/uwsgi
mkdir -p /tmp/nginx/scgi

bashio::log.info "Starting nginx with multi-location asset handling strategy..."
bashio::log.info "Using separate location blocks for CSS, JS, and image assets"
bashio::log.info "Forcing correct MIME types to prevent browser rejection"
bashio::log.info "This should resolve 'MIME type not supported' errors"

# Add debugging to nginx config
echo "Starting nginx with debug-level error logging..." >&2

# Start nginx
exec nginx -c /etc/nginx/nginx.conf