# kiwix/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Kiwix
# Runs nginx for ingress support
# ==============================================================================

# Replace ingress entry placeholder
sed -i "s|%%ingress_entry%%|$(bashio::addon.ingress_entry)|g" /etc/nginx/nginx.conf

# Wait for Kiwix to be available
while ! nc -z 127.0.0.1 8090; do
    bashio::log.info "Waiting for Kiwix to start..."
    sleep 2
done

bashio::log.info "Starting nginx for ingress support..."

# Test nginx configuration
if ! nginx -t; then
    bashio::log.error "nginx configuration test failed!"
    exit 1
fi

# Start nginx
exec nginx
